# system/crawl1.service
[Unit]
After=network.service
Description=Crawling @ Home

[Service]
Type=simple
LimitNOFILE=65535
WorkingDirectory=/home/crawl
ExecStart=/home/crawl/crawl1.sh
User=crawl

[Install]
WantedBy=multi-user.target
# WantedBy=default.target